# supervisor/web.conf
[program:web]
command=python3 -m http.server 10000
directory=/home/%(ENV_SSH_USER)s
autostart=true
autorestart=true
